# foreign-server-nginx/crontab/cron_job.sh
# m h  dom mon dow   command
0 5  * * *  /srv/nginx/cron_job.sh
