# domestic-server-v2ray/crontab
# m    h dom mon dow   command
  */5  *  *   *   *    /srv/v2fly/v2stat.py
